# conf/supervisord_cbe.conf
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".

[unix_http_server]
file=/home/cbe-master/realfast/soft/supervisor.sock  ; (the path to the socket file)

[supervisord]
;logfile=/home/cbe-master/realfast/soft/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/home/cbe-master/realfast/soft/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
childlogdir=/home/cbe-master/realfast/soft
strip_ansi=true

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/cbe-master/realfast/soft/supervisor.sock  ; use a unix:// URL  for a unix socket

;;;;;;;;;;;;
; programs ;
;;;;;;;;;;;;

[program:redis]
command=env PATH="/home/cbe-master/realfast/anaconda/envs/deployment/bin:%(ENV_PATH)s" redis-server /home/cbe-master/realfast/soft/realfast/conf/redis_cbe.conf
process=%(program_name)s
priority=1    ; set up db first

redirect_stderr=True
autostart=False

[program:mcaf]
; With the below command (includes -N), NRAO controls archiving
command=env PATH="/home/cbe-master/realfast/anaconda/envs/deployment/bin:%(ENV_PATH)s" nice -n10 mcaf_monitor -v --production --intent "TARGET" --project "17B-301" --rtparams /lustre/evla/test/realfast/rtpipe_cbe2.conf 
;command=nice -n10 mcaf_monitor -v --production --intent "TARGET" --project "15B-305" --rtparams /lustre/evla/test/realfast/rtpipe_cbe.conf
;command=nice -n10 mcaf_monitor -v --slow 5 --production --intent "" --project "" ; for tests
process=%(program_name)s
priority=3     ; watch for data last
;environment=PATH=/home/cbe-master/realfast/anaconda/envs/deployment/bin:%(ENV_PATH)s  ; use deployment environment
redirect_stderr=True
stdout_logfile=/home/cbe-master/realfast/soft/realfast.log       ; stdout log path, NONE for none; default AUTO
stdout_events_enabled=true   ; emit events on stderr writes (default false)
buffer_size = 1024
autostart=False

[program:queue]
;command=queue_monitor -a -v --triggered --threshold 7.0 --production --bdfdir /lustre/evla/wcbe/data/no_archive
; Another "NRAO controls archiving"
command=env PATH="/home/cbe-master/realfast/anaconda/envs/deployment/bin:%(ENV_PATH)s" queue_monitor -a -v --threshold 7.5 --production --bdfdir /lustre/evla/wcbe/data/no_archive
process=%(program_name)s
priority=2
;environment=PATH=/home/cbe-master/realfast/anaconda/envs/deployment/bin  ; use deployment environment
redirect_stderr=True
stdout_logfile=/home/cbe-master/realfast/soft/realfast.log       ; stdout log path, NONE for none; default AUTO
stdout_events_enabled=true   ; emit events on stderr writes (default false)
buffer_size = 1024
autostart=False

[program:work2]
command=ssh -t -t cbe-node-02 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
;command=ssh -t -t cbe-node-02 'source activate deployment; python -c "import sdmpy; print(sdmpy.__file__)"'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work3]
command=ssh -t -t cbe-node-03 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work4]
command=ssh -t -t cbe-node-04 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work5]
command=ssh -t -t cbe-node-05 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work6]
command=ssh -t -t cbe-node-06 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work7]
command=ssh -t -t cbe-node-07 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work8]
command=ssh -t -t cbe-node-08 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work9]
command=ssh -t -t cbe-node-09 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work10]
command=ssh -t -t cbe-node-10 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work11]
command=ssh -t -t cbe-node-11 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work12]
command=ssh -t -t cbe-node-12 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work13]
command=ssh -t -t cbe-node-13 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work14]
command=ssh -t -t cbe-node-14 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work15]
command=ssh -t -t cbe-node-15 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work16]
command=ssh -t -t cbe-node-16 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work17]
command=ssh -t -t cbe-node-17 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work18]
command=ssh -t -t cbe-node-18 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work19]
command=ssh -t -t cbe-node-19 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work20]
command=ssh -t -t cbe-node-20 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work21]
command=ssh -t -t cbe-node-21 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work22]
command=ssh -t -t cbe-node-22 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work23]
command=ssh -t -t cbe-node-23 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:work28]
command=ssh -t -t cbe-node-28 'source activate deployment; nice -n10 rq worker default -u redis://cbe-node-01'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

; shortcut for all workres. controlled, e.g., 'supervisorctl start workers:*'
[group:workers]
;programs=work2,work3,work4,work5,work6,work7,work8,work9,work10,work11,work12,work13,work14,work15,work16,work17,work18,work19 ;,work20,work21 ;,work22 ;,work23
programs=work5,work6,work7,work8,work9,work10,work11,work12,work13,work14,work15,work16,work17,work18,work19 ;,work20,work21 ;,work22 ;,work23

;;;;;;;;;;;;;;
;; realfast ;;
;;;;;;;;;;;;;;

[program:dsched]
command=ssh -t -t cbe-node-01 'source activate development; dask-scheduler'
;command=ssh -t -t cbe-node-01 'source activate development; dask-scheduler --interface ib0'
process=%(program_name)s
environment=PATH=/home/cbe-master/realfast/anaconda/envs/development/bin:%(ENV_PATH)s  ; use deployment environment
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:rf001gpu1]
command=ssh -t -t rfnode001 'source activate development; dask-worker 192.168.201.101:8786 --interface ib0 --name rf001gpu1 --nprocs 1 --nthreads 18 --memory-limit=22e9 --local-directory /lustre/evla/test/realfast/scratch --resources "GPU=1,MEMORY=32,CORES=18"'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:rf001gpu2]
command=ssh -t -t rfnode001 'source activate development; dask-worker 192.168.201.101:8786 --interface ib0 --name rf001gpu2 --nprocs 1 --nthreads 18 --memory-limit=22e9 --local-directory /lustre/evla/test/realfast/scratch --resources "GPU=1,MEMORY=32,CORES=18"'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:rf001r1]
command=ssh -t -t rfnode001 'source activate development; dask-worker 192.168.201.101:8786 --interface ib0 --nprocs 1 --nthreads 1 --memory-limit=10e9 --local-directory /lustre/evla/test/realfast/scratch --resources "READER=1"'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:rf001r2]
command=ssh -t -t rfnode001 'source activate development; dask-worker 192.168.201.101:8786 --interface ib0 --nprocs 1 --nthreads 1 --memory-limit=10e9 --local-directory /lustre/evla/test/realfast/scratch --resources "READER=1"'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:cbe02]
command=ssh -t -t cbe-node-02 'source activate development; dask-worker 192.168.201.101:8786 --interface ib0 --nprocs 1 --nthreads 15 --memory-limit=14e9 --local-directory /lustre/evla/test/realfast/scratch --resources "MEMORY=24,CORES=15"'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:cbe02r1]
command=ssh -t -t rfnode001 'source activate development; dask-worker 192.168.201.101:8786 --interface ib0 --nprocs 1 --nthreads 1 --memory-limit=10e9 --local-directory /lustre/evla/test/realfast/scratch --resources "READER=1"'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:cbe03]
command=ssh -t -t cbe-node-03 'source activate development; dask-worker 192.168.201.101:8786 --interface ib0 --nprocs 1 --nthreads 15 --memory-limit=14e9 --local-directory /lustre/evla/test/realfast/scratch --resources "MEMORY=24,CORES=15"'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:cbe03r1]
command=ssh -t -t rfnode001 'source activate development; dask-worker 192.168.201.101:8786 --interface ib0 --nprocs 1 --nthreads 1 --memory-limit=10e9 --local-directory /lustre/evla/test/realfast/scratch --resources "READER=1"'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:cbe04]
command=ssh -t -t cbe-node-04 'source activate development; dask-worker 192.168.201.101:8786 --interface ib0 --nprocs 1 --nthreads 15 --memory-limit=14e9 --local-directory /lustre/evla/test/realfast/scratch --resources "MEMORY=24,CORES=15"'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[program:cbe04r1]
command=ssh -t -t rfnode001 'source activate development; dask-worker 192.168.201.101:8786 --interface ib0 --nprocs 1 --nthreads 1 --memory-limit=10e9 --local-directory /lustre/evla/test/realfast/scratch --resources "READER=1"'
process=%(program_name)s
priority=5     ; start worker after db
redirect_stderr=True
autostart=False

[group:rf]
programs=dsched,rf001gpu1,rf001gpu2,rf001r1,rf001r2 ;,cbe02,cbe02r1,cbe03,cbe03r1,cbe04,cbe04r1

[program:config_catcher]
command=realfast config_catcher --preffile /home/cbe-master/realfast/lustre_workdir/realfast.yml
redirect_stderr=True
environment=PATH="/home/cbe-master/realfast/anaconda/envs/development/bin:${PATH}"  ; use deployment environment
autostart=False
